# overlay/etc/service/apache2/run
#!/bin/sh
exec /usr/bin/apache2 -DFOREGROUND >>/var/log/apache.log 2>&1
