# overlay/etc/service/postgres/run
#!/bin/sh
exec /usr/lib/postgresql/9.3/bin/postgres -D /var/lib/postgresql/9.3/main -c config_file=/etc/postgresql/9.3/main/postgresql.conf >> /var/log/postgresql.log 2>&1
