# overlay/etc/service/lcr/run
#!/bin/sh
exec /usr/sbin/lcr start >> /var/log/lcr.log 2>&1
